# mapping_mediator/crontab
0 22 * * * root python /app/apps/srilanka/view/export_ocl_cron.py >> /var/log/cron.log 2>&1
